# rootfs/etc/s6-overlay/s6-rc.d/vdlm_server/run
#!/bin/sh
exec /etc/s6-overlay/scripts/vdlm_server
